# config/systemd/github-runner@.service
[Unit]
Description=GitHub Actions Runner %i
After=network.target

[Service]
User=github-runner
Group=github-runner
WorkingDirectory=/opt/actions-runner/%i
ExecStart=/opt/actions-runner/%i/runsvc.sh
KillMode=process
Restart=always
RestartSec=5

# Security hardening (baseline)
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectSystem=strict
ProtectHome=true
ProtectControlGroups=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
LockPersonality=yes
RestrictSUIDSGID=yes
SystemCallFilter=@system-service
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
CapabilityBoundingSet=
ReadWritePaths=/opt/actions-runner/%i
ReadWritePaths=/var/lib/github-runner/%i

[Install]
WantedBy=multi-user.target
