# config/systemd/github-runner@.service
[Unit]
Description=GitHub Actions Runner %i
Wants=network-online.target
After=network-online.target

[Service]
User=github-runner
Group=github-runner
WorkingDirectory=/opt/actions-runner/%i

# CRITICAL: Use exec without -l flag to prevent environment override
ExecStart=/bin/bash -c 'cd /opt/actions-runner/%i && exec ./run.sh --startuptype service'

KillMode=process
Restart=always
RestartSec=5

# Baseline hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=yes
ProtectControlGroups=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
LockPersonality=yes
RestrictSUIDSGID=yes
CapabilityBoundingSet=
AmbientCapabilities=

[Install]
WantedBy=multi-user.target